# service/svscan/run
#!/bin/sh

/usr/sbin/openmqtt -c /etc/openmqtt/openmqtt.conf
