# service/svscan/run
#!/bin/sh

/usr/sbin/minion -c /etc/minion/minion.conf
