# service/upstart/minion.conf
description "Minion v0.0.1"
author "Moran, Inc. <daniel@moranit.com>"

start on net-device-up

respawn

exec /usr/sbin/minion -c /etc/minion/minion.conf
