# service/upstart/openmqtt.conf
description "OpenMqtt v0.0.1"
author "OpenMqtt <daniel@openmqtt.org>"

start on net-device-up

respawn

exec /usr/sbin/openmqtt -c /etc/openmqtt/openmqtt.conf
