# service/upstart/openmqtt.conf
description "OpenMqtt v0.0.1"
author "Moran, Inc. <daniel@moranit.com>"

start on net-device-up

respawn

exec /usr/sbin/openmqtt -c /etc/openmqtt/openmqtt.conf
